# db/supervisord.conf
[supervisord]
nodaemon=true


[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true


[program:db_start]
command=/home/coco/start.sh
autostart=true
user=coco
priority=1

